# s6/syslogd/run
#!/usr/bin/env sh

exec syslogd -n -O /dev/stdout -S
