# s6/xinetd/run
#!/usr/bin/env sh

exec s6-applyuidgid -u 1000 -g 1000 xinetd -dontfork
